# io.cemantix.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.cemantix.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/johan/Documents/code/cemantix/run_daily.sh</string>
  </array>

  <!-- Lancement chaque jour à 00h15 (heure locale) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <!-- Si le Mac était éteint à l'heure prévue, lancer au prochain démarrage -->
  <key>RunAtLoad</key>
  <false/>

</dict>
</plist>
